# debian/postinst
#!/bin/bash
update-rc.d -f rpimonitor defaults
/etc/init.d/rpimonitor start
