# debian/prerm
#!/bin/bash
if [ "$1" = "remove" ]; then
	/etc/init.d/rpimonitor stop
	update-rc.d -f rpimonitor remove
fi
